# deploy/reasonix-telegram.service
[Unit]
Description=reasonix-telegram: Telegram bridge for Reasonix
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root
EnvironmentFile=/etc/reasonix-telegram.env
ExecStart=/usr/local/bin/reasonix-telegram
Restart=on-failure
RestartSec=5s
TimeoutStopSec=310

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
