# deploy/reasonix-telegram.service
[Unit]
Description=reasonix-telegram: Telegram bridge for Reasonix
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root
EnvironmentFile=/etc/reasonix-telegram.env
Environment=REASONIX_CACHE_DIR=/var/lib/reasonix-telegram/cache
Environment=XDG_CACHE_HOME=/var/lib/reasonix-telegram/cache
ExecStart=/usr/local/bin/reasonix-telegram
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15

NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/var/lib/reasonix-telegram
ProtectHome=read-only
PrivateTmp=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target